# kismet-base/config/skel/.config/autostart/kismet-extra-polish.desktop
[Desktop Entry]
Type=Application
Name=Kismet Extra Polish
Comment=Apply extra Kismet desktop polish on first login
Exec=/usr/local/bin/kismet-extra-polish
Terminal=false
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
