# kismet-base/config/skel/.config/autostart/kismet-plasma-polish.desktop
[Desktop Entry]
Type=Application
Version=1.0
Name=Kismet Plasma Polish
Comment=Apply Kismet desktop defaults on first login
Exec=/usr/local/bin/kismet-plasma-polish
OnlyShowIn=KDE;
X-KDE-autostart-phase=2
NoDisplay=true
