# kismet-base/config/skel/.config/autostart/kismet-gnome-polish.desktop
[Desktop Entry]
Type=Application
Name=Kismet GNOME Polish
Comment=Apply Zorin-style GNOME theming
Exec=kismet-gnome-polish
Terminal=false
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true