# kismet-base/config/etc/systemd/system/kismet-agent.service
[Unit]
Description=Kismet OS preview system agent
After=network-online.target ollama.service kismet-ollama-loader.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/lib/kismet/kismet_agent.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/usr/local/lib/kismet

[Install]
WantedBy=multi-user.target
